# audio-detect-daemon.service
[Unit]
Description=Audio Stream Detection Daemon
Documentation=file:///home/cgeng/work/myfiles/audio-capture/docs/daemon-installation.md
After=pulseaudio.service sound.target
Wants=pulseaudio.service

[Service]
Type=simple
User=%i

# Environment for audio access
Environment="PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/%i"
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="PULSE_RUNTIME_PATH=/run/user/%U/pulse"
Environment="PULSE_SERVER=unix:/run/user/%U/pulse/native"
Environment="DISPLAY=:0"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"

# Working directory
WorkingDirectory=%h/.local/share/audio-detect

# Command
ExecStart=/usr/bin/env -i PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin HOME=/home/%i XDG_RUNTIME_DIR=/run/user/%U PULSE_RUNTIME_PATH=/run/user/%U/pulse PULSE_SERVER=unix:/run/user/%U/pulse/native DISPLAY=:0 DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus /home/%i/.local/bin/uvx --from git+https://github.com/ChristianGeng/audio-capture.git audio-detect-daemon

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=audio-detect-daemon

[Install]
WantedBy=default.target
