# teams-auto-recorder.service
[Unit]
Description=Automatic Teams Meeting Audio Recorder
Documentation=file:///home/cgeng/work/myfiles/audio-capture/docs/adr/0001-simplify-teams-audio-capture.md
After=pulseaudio.service sound.target
Wants=pulseaudio.service

[Service]
Type=simple
User=cgeng
Group=cgeng

# Environment
Environment="PATH=/home/cgeng/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HOME=/home/cgeng"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_RUNTIME_PATH=/run/user/1000/pulse"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"

# Working directory
WorkingDirectory=/home/cgeng/work/myfiles/audio-capture

# Command
ExecStart=/usr/bin/python3 /home/cgeng/work/myfiles/audio-capture/teams_auto_recorder.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=teams-recorder

[Install]
WantedBy=default.target
